# run_domino.service
[Unit]
Description=for my robot
After=network-online.target

[Service]
ExecStart=/home/domino/robocon/dog.sh

[Install]
WantedBy=multi-user.target